# cron/crossposting.cron
*/5 * * * * root curl --silent http://pipeinpipe.info/cron/cron_crossposting.php?password=yreutywerUIUihkhsdjahsduhweuUUUUU
